# files/rss-notifier.service
[Unit]
Description=RSS Notifier Service

[Service]
ExecStart=/opt/rss-notifier.sh
