# files/rss-notifier.timer
[Unit]
Description=RSS Notifier Timer

[Timer]
OnCalendar=*-*-* *:*:0
Persistent=true

[Install]
WantedBy=timers.target